# contrib/energy-manager.service
[Unit]
Description=Energy Manager — Gestionnaire d'énergie Rust (remplace Node-RED)
Documentation=https://github.com/thieryus007-cloud/Daly-BMS-Rust
After=network-online.target mosquitto.service
Wants=network-online.target
PartOf=daly-bms.service

[Service]
Type=simple
User=pi5compute
Group=pi5compute

ExecStart=/usr/local/bin/energy-manager
WorkingDirectory=/home/pi5compute/Daly-BMS-Rust

# Config
Environment=RUST_LOG=info
EnvironmentFile=-/etc/daly-bms/.env
Environment=ENERGY_CONFIG=/etc/daly-bms/config.toml

Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
MemoryMax=100M

[Install]
WantedBy=multi-user.target
